# linux/tasky-api.service
[Unit]
Description=Tasky API

[Service]
ExecStart=~/tasky/tasky-api

[Install]
WantedBy=multi-user.target 
