# linux/tasky.service

[Unit]
Description=Tasky

[Service]
ExecStart=~/tasky/tasky-linux-x64/tasky

[Install]
WantedBy=multi-user.target 
